# etc/s6/services/01-django/run
#!/bin/sh
NAME="DJANGO" 
COMMAND="cd /data && python3 manage.py runserver 0.0.0.0:9000" 
WAIT="/tmp/django"


exec service-handler "$NAME" "$COMMAND" "$WAIT"
